# deploy/spiregithubtoken.timer
[Unit]
Description=Run Spire GitHub token generator every 2 minutes

[Timer]
OnBootSec=10s
OnUnitActiveSec=2m

[Install]
WantedBy=timers.target
